# root/etc/services.d/autorippr/run
#!/usr/bin/with-contenv bash

echo "Starting Autorippr Disc Ripper."
exec s6-setuidgid root /bin/bash -c \
 'python /Autorippr-1.7.0/autorippr.py --all'
