# services.d/autorippr/run
#!/usr/bin/env
exec s6-setuidgid abc python /Autorippr-1.7.0/autorippr.py --all
