# rootfs/etc/services.d/autorippr/run
#!/usr/bin/with-contenv bash

echo "Starting Autorippr Disc Ripper."
exec python /Autorippr-1.7.0/autorippr.py --all
